# systemd/openclaw-healthcheck.service
[Unit]
Description=OpenClaw self-healing health check
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/openclaw-healthcheck.sh
StandardOutput=journal
StandardError=journal
